# render/fleet/test-fixtures/units/instance-count-1/simple.hcl/test-web-backup-ti@2.timer
[Unit]
Description=Timer unit for test/web/backup [slice 2]

[Timer]
OnCalendar=hourly
Unit=test-web-backup-mn@2.service

[X-Fleet]
MachineOf=test-web-backup-mn@2.service

[X-testproject]
GeneratedBy="testproject test-version, build test-build"
